# quillforge-watcher.path
[Unit]
Description=Watch QuillForge binary for changes

[Path]
PathChanged=/mnt/den-data/dev/quillforge/release/QuillForge.Web
Unit=quillforge-restart.service

[Install]
WantedBy=multi-user.target
